# supervisord.conf
[supervisord]
nodaemon=true

[program:sbbs]
directory=/opt/schemebbs
command=scheme --args 8080 < bbs.scm

[program:nginx]
command=nginx -g 'daemon off;'
