# supervisord.conf
[supervisord]
nodaemon=true

[program:sbbs]
directory=/opt/schemebbs
command=/usr/local/bin/scheme --args 8080 < bbs.scm

[program:nginx]
command=/usr/local/nginx/sbin/nginx -g 'daemon off;'
